# kmsg-backend.service
[Unit]
Description=KMSG Backend Service
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ubuntu/kmsg/backend
ExecStart=/bin/bash -c 'cd /home/ubuntu/kmsg/backend && docker-compose up -d'
ExecStop=/bin/bash -c 'cd /home/ubuntu/kmsg/backend && docker-compose down'
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
